# Gen-code/nextest_cgen-r4-0a-pre47/iServer/utils/cdrpp.cron
30 1 * * *  	(cd /usr/local/nextone/etc ; ./cdrpp -i cin -o outf -n mail.pl)
